# extras/supervisord/supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ;
logfile_maxbytes=50MB        ;
logfile_backups=10           ;
loglevel=info                ;
pidfile=/tmp/supervisord.pid ;
nodaemon=false               ;
minfds=1024                  ;
minprocs=200                 ;


[program:gunicorn]
command=/simple-route/walmart/bin/gunicorn app.wsgi:application --bind 127.0.0.1:7000 --pid /tmp/app.pid ;
directory=/simple-route/app/ ;


[inet_http_server]
port=127.0.0.1:9001   ;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001 ;
